# lava-lights
### BEGIN INIT INFO
# Provides: lava-lights
# Required-Start:
# Required-Stop:
# Default-Start:        2 3 4 5
# Default-Stop:         0 1 6
# Short-Description:    Lavalamp build listener
### END INIT INFO

export PATH=$PATH:/usr/local/bin/node
export NODE_PATH=$NODE_PATH:/usr/local/lib/node_modules
export HOME=/root

case "$1" in
  start)
    /usr/local/bin/forever start --workingDir /home/pi/lava-lights --sourceDir /home/pi/lava-lights -a -l forever.log -o out.log -e err.log --minUptime 1000 --spinSleepTime 1000 -c "node --harmony_destructuring --harmony_rest_parameters" index.js
    ;;
  stop)
    exec /usr/local/bin/forever stopall
    ;;
  *)

  echo "Usage: /etc/init.d/nodeup {start|stop}"
  exit 1
  ;;
esac
exit 0